# meta-tegra-extensions/recipes-bsp/tegra-rootfs-redundancy/files/wendyos-tegra-rootfs-redundancy.service
[Unit]
Description=Arm Jetson rootfs A/B redundancy if firmware left it disabled
DefaultDependencies=no
# Needs the persistent data partition for its one-shot attempt marker, and
# efivarfs for the UEFI variable. Must settle the slot machinery BEFORE the
# OTA boot verifier reads/commits a pending update, so it is ordered ahead of
# wendyos-update-verify.service. It may reboot once on first run; that is
# guarded so it can never loop.
After=local-fs.target data.mount
# The attempt marker lives on /data; RequiresMountsFor pulls the mount into the
# transaction and fails this unit if /data is unavailable, so the marker can
# never land on the rootfs underlay and get hidden once /data mounts (which
# would defeat the reboot-loop guard). After= alone does not guarantee this.
RequiresMountsFor=/data
Before=wendyos-update-verify.service
ConditionPathExists=/sys/firmware/efi/efivars

[Service]
Type=oneshot
ExecStart=/usr/sbin/wendyos-tegra-arm-rootfs-redundancy
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
