# meta-x86-extensions/recipes-core/wendyos-identity/files/wendyos-device-name-generate.service
[Unit]
Description=WendyOS Device Name Generation
# x86 variant of the shared wendyos-identity unit. Like the base (Tegra/RPi)
# unit but without the wendyos-etc-binds / /data ordering, which does not exist
# on x86 (WENDYOS_OTA=none) — identity lands on the persistent rootfs. Still
# ordered after the UUID generator (the name derives from the UUID) and before
# identity publishes. Keep [Service]/[Install] in sync with
# recipes-core/wendyos-identity/files/wendyos-device-name-generate.service.
After=wendyos-uuid-generate.service
Before=wendyos-identity.service

[Service]
Type=oneshot
ExecStart=/usr/bin/generate-device-name.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

