# recipes-core/systemd-mount-wendy/files/var-lib-wendy.mount
[Unit]
Description=Bind mount /var/lib/wendy (persistent app volumes) from data partition
Documentation=man:systemd.mount(5)

# The agent creates persistent app volumes under /var/lib/wendy/volumes
# (see wendy-agent oci/entitlements.go). By default that path is on the small
# A/B rootfs, so heavy volume data fills the rootfs and is wiped on OTA. Back
# it with /data/wendy (persists across A/B switches), same pattern as
# /var/lib/containerd -> /data/containerd and /var/lib/tee -> /data/tee.

# Ensure data partition is mounted and filesystem is fully expanded
RequiresMountsFor=/data
After=data.mount

# Wait for the /data fs grow so first-boot volume data can't hit ENOSPC while the
# fs is still small. Growers are board-specific, inert where absent:
# mender-systemd-growfs-data (Tegra), grow-data-fs-online (RPi).
After=mender-systemd-growfs-data.service grow-data-fs-online.service

# Must be mounted before the agent starts
Before=wendyos-agent.service
DefaultDependencies=no

[Mount]
What=/data/wendy
Where=/var/lib/wendy
Type=none
Options=bind,nosuid,nodev,x-systemd.mkdir

[Install]
WantedBy=local-fs.target
