# meta-x86-extensions/recipes-support/wendyos-driver-installer/files/wendyos-driver-update.service
[Unit]
Description=WendyOS feed-backed GPU driver update
Wants=network-online.target
After=network-online.target NetworkManager-wait-online.service
ConditionPathExists=/etc/wendyos/driver-request.conf

[Service]
Type=oneshot
ExecStart=/usr/sbin/wendyos-driver-update.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
