# recipes-core/systemd-mount-tee/files/var-lib-tee.mount
[Unit]
Description=Bind mount /var/lib/tee (OP-TEE secure storage) from persistent data partition
Documentation=man:systemd.mount(5)

# OP-TEE secure storage (PKCS#11 tokens, device keys, certificates) is written
# by tee-supplicant to /var/lib/tee. By default that path is on the rootfs, which
# an A/B OTA replaces on every update -> the device's keys would be wiped.
# Back it with /data/tee (persists across A/B switches), same pattern as
# /var/lib/containerd -> /data/containerd.

# Ensure data partition is mounted and filesystem is fully expanded
RequiresMountsFor=/data
After=data.mount


# Must be mounted before tee-supplicant uses /var/lib/tee for secure storage.
# Assumes source-built OP-TEE (USE_PREBUILT_OPTEE=0, the meta-tegra default and
# this project's config): provider optee-client, service tee-supplicant.service,
# CFG_TEE_FS_PARENT_PATH=/var/lib/tee. The prebuilt provider (nvidia-l4t-optee,
# USE_PREBUILT_OPTEE=1) uses nv-tee-supplicant.service and an NVIDIA-controlled
# storage path — we order before it too so the ordering is not silently lost if
# the provider is ever switched, but the /var/lib/tee path below would then need
# re-checking.
Before=tee-supplicant.service nv-tee-supplicant.service
DefaultDependencies=no

[Mount]
What=/data/tee
Where=/var/lib/tee
Type=none
Options=bind,nosuid,nodev,x-systemd.mkdir

[Install]
WantedBy=local-fs.target
