# recipes-core/wendyos-identity/files/wendyos-device-name-generate.service
[Unit]
Description=WendyOS Device Name Generation
# Same ordering rule as wendyos-uuid-generate: run after the /data bind is in
# place so device-name lands on the persistent /etc/wendyos, and before identity
# consumes it. Was DefaultDependencies=no + WantedBy=sysinit.target, which
# front-ran the bind (see wendyos-uuid-generate.service).
After=wendyos-etc-binds.service wendyos-uuid-generate.service
Requires=wendyos-etc-binds.service
Before=wendyos-identity.service

[Service]
Type=oneshot
ExecStart=/usr/bin/generate-device-name.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

