# recipes-core/systemd-mount-home/files/home.mount
[Unit]
Description=Bind mount /home from persistent data partition
Documentation=man:systemd.mount(5)
# Ensure data partition is mounted and filesystem is fully expanded
RequiresMountsFor=/data
After=data.mount
# Wait for Mender data partition resize to complete (if enabled)
# This is safe even if mender-systemd-growfs-data.service doesn't exist
After=mender-systemd-growfs-data.service
Before=local-fs.target

[Mount]
What=/data/home
Where=/home
Type=none
# x-systemd.mkdir creates the /data/home bind source before mounting. On a
# freshly initialised /data (first boot, or after data-init reinitialises a
# bad fs) the source does not exist yet, and a bind mount of a missing source
# fails — leaving /home as the empty rootfs skeleton and breaking the wendy
# user's persistent home (and, via linger, user@1000.service). Every other
# /data bind mount (var-lib-containerd, var-lib-tee, var-lib-wendy, var-log)
# already uses this; home.mount was the lone exception.
Options=bind,x-systemd.mkdir

[Install]
WantedBy=local-fs.target
