# recipes-core/wendyos-data-setup/files/data.mount
[Unit]
Description=WendyOS persistent data partition (/data)
# Format/grow on first boot must complete before we mount.
Requires=wendyos-data-init.service
After=wendyos-data-init.service

[Mount]
What=/dev/disk/by-partlabel/data
Where=/data
Type=ext4
# Hardened (H6): nosuid,nodev. NOT noexec — /data backs container images and app
# volumes (bind-mounted to /var/lib/containerd, /var/lib/wendy, /home) that execute.
Options=defaults,noatime,nosuid,nodev

[Install]
WantedBy=local-fs.target
