# meta-rpi-extensions/recipes-support/grow-data-part/files/grow-data-part.service
[Unit]
Description=WendyOS grow /data PARTITION to fill the storage device (offline, first boot)
# Phase 1: grow the /data PARTITION while it's still unmounted. Ordering before
# local-fs-pre.target puts us ahead of every mount and systemd-fsck, so parted has
# no race. The slow ext4 grow is split into grow-data-fs-online.service (off the
# boot path). Idempotent: parted no-ops once /data fills the disk.
DefaultDependencies=no
Before=local-fs-pre.target data.mount
Conflicts=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Fast phase (parted + clean-flag fsck); a tight bound fails a wedged tool rather
# than stalling boot.
TimeoutStartSec=120
ExecStart=/usr/sbin/grow-data-part.sh partition

[Install]
WantedBy=data.mount
